# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
nodaemon=true
[program:ondemand]
command=/opt/rh/httpd24/root/usr/sbin/httpd-scl-wrapper -D FOREGROUND
[program:keycloak]
command=/opt/keycloak-4.8.3.Final/bin/standalone.sh -b 0.0.0.0