# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
nodaemon=true
[program:ondemand]
command=/opt/rh/httpd24/root/usr/sbin/httpd-scl-wrapper -D FOREGROUND